# board/speedsaver/rootfs_overlay/etc/systemd/system/navit.service
[Unit]
Description=navit software

[Service]
ExecStart=/usr/bin/navit
Restart=on-failure

[Install]
WantedBy=multi-user.target
